# supervisor/conf.d/evolution-api.conf
[program:evolution-api]
command=/bin/bash -c "./Docker/scripts/deploy_database.sh && npm run start:prod"
autorestart=false
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
