# supervisor/conf.d/run.conf
[program:run]
command=/bin/bash -c '/evolution/run.sh 43200'
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
